# daemons/autocopter.service
[Unit]
Description=Autocopter

[Service]
ExecStart=/home/pi/autocopter/autocopter.py
Restart=on-abort

[Install]
WantedBy=multi-user.target
